# scripts/launchd/com.farfield.app.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.farfield.app</string>

    <key>ProgramArguments</key>
    <array>
      <string>__HOME__/bin/run-farfield-app.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__FARFIELD_DIR__</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/farfield/app.out.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/farfield/app.err.log</string>

    <key>ProcessType</key>
    <string>Background</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
  </dict>
</plist>
